# contrib/systemd/bdstorage-dedupe-once.service
# Triggered by bdstorage-dedupe.timer (or: systemctl start bdstorage-dedupe-once.service).

[Unit]
Description=bdstorage one-shot dedupe

[Service]
Type=oneshot
User=bdstorage
Group=bdstorage
EnvironmentFile=-/etc/bdstorage.env
ExecStart=/usr/local/bin/bdstorage dedupe $BDSTORAGE_TARGET
