# contrib/wsl/bdstorage-dedupe.timer
# User timer: runs bdstorage-dedupe-once.service on a schedule (alternative to the long-running daemon).

[Unit]
Description=Timer for bdstorage dedupe (WSL user)

[Timer]
OnCalendar=daily
Unit=bdstorage-dedupe-once.service
Persistent=true

[Install]
WantedBy=timers.target
