# contrib/wsl/bdstorage-dedupe.service
# User unit for WSL / Ubuntu: install with contrib/wsl/install-wsl-user-units.sh
#
# State DB and vault follow HOME (see src/state.rs). Override in ~/.config/bdstorage/env if needed.

[Unit]
Description=bdstorage deduplication daemon (WSL user)
After=default.target

[Service]
Type=simple
EnvironmentFile=-%h/.config/bdstorage/env
ExecStart=$BDSTORAGE_BIN daemon $BDSTORAGE_TARGET --interval-secs $BDSTORAGE_INTERVAL_SECS
Restart=on-failure
RestartSec=30

[Install]
WantedBy=default.target
