# contrib/wsl/bdstorage-dedupe-once.service
# User oneshot: triggered by bdstorage-dedupe.timer (no long-running daemon).

[Unit]
Description=bdstorage one-shot dedupe (WSL user)
After=default.target

[Service]
Type=oneshot
EnvironmentFile=-%h/.config/bdstorage/env
ExecStart=$BDSTORAGE_BIN dedupe $BDSTORAGE_TARGET
